# graphics/mesa-libs/pkg-plist
etc/libmap.d/mesa.conf
include/EGL/eglext_angle.h
include/EGL/eglmesaext.h
@comment include/GL/internal/dri_interface.h
include/gbm.h
include/gbm_backend_abi.h
@comment lib/dri/crocus_dri.so
@comment lib/dri/i915_dri.so
@comment lib/dri/iris_dri.so
@comment lib/dri/libdril_dri.so
@comment lib/dri/r300_dri.so
@comment lib/dri/r600_dri.so
@comment lib/dri/r600_drv_video.so
@comment lib/dri/radeonsi_dri.so
@comment lib/dri/radeonsi_drv_video.so
@comment lib/dri/virtio_gpu_dri.so
@comment lib/dri/virtio_gpu_drv_video.so
@comment lib/dri/vmwgfx_dri.so
@comment lib/dri/zink_dri.so
lib/gbm/dri_gbm.so
lib/libEGL_mesa.so
lib/libEGL_mesa.so.0
lib/libEGL_mesa.so.0.0.0
%%GLX%%lib/libGLX_mesa.so
%%GLX%%lib/libGLX_mesa.so.0
%%GLX%%lib/libGLX_mesa.so.0.0.0
lib/libgallium-25.2.1.so
lib/libgbm.so
lib/libgbm.so.1
lib/libgbm.so.1.0.0
@comment libdata/pkgconfig/dri.pc
libdata/pkgconfig/gbm.pc
@comment share/drirc.d/00-mesa-defaults.conf
share/glvnd/egl_vendor.d/50_mesa.json
